# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:squeezeboxserver]
command=squeezeboxserver --user abc --prefsdir /config/prefs --logdir /config/logs --cachedir /config/cache

